# supervisord.conf
[supervisord]
nodaemon = true

[program:minecraft]
command = java -Xmx4096M -Xms1024M -jar /minecraft/minecraft-server.jar nogui
directory = /world
